# soc-bot/soc-bot.service
[Unit]
Description=SOC Bot (Wazuh Alert Processor)
After=network.target

[Service]
Type=simple
User=__SOC_BOT_USER__
WorkingDirectory=__SOC_BOT_DIR__
EnvironmentFile=__SOC_BOT_DIR__/.env
ExecStart=__SOC_BOT_DIR__/venv/bin/python main.py
Restart=always
RestartSec=10
StartLimitIntervalSec=60
StartLimitBurst=3
StandardOutput=journal
StandardError=journal
SyslogIdentifier=soc-bot

[Install]
WantedBy=multi-user.target
